# etc/init.d/rsyslog-openrc
#!/sbin/openrc-run

description="Rsyslog is an enhanced multi-threaded syslogd"
description_start="Starting enhanced syslogd"
description_stop="Stopping enhanced syslogd"

DAEMON=/usr/sbin/rsyslogd
PIDFILE=/run/rsyslogd.pid
RSYSLOGD_OPTIONS=""

depend() {
    need mountall
    after bootmisc
}

sendsigs_omit() {
    OMITDIR=/run/sendsigs.omit.d
    mkdir -p $OMITDIR
    ln -sf $PIDFILE $OMITDIR/rsyslog
}

start() {
    ebegin "${description_start}"
    start-stop-daemon --start --quiet --pidfile ${PIDFILE} --exec ${DAEMON} -- ${RSYSLOGD_OPTIONS}
    sendsigs_omit
    eend $?
}

stop() {
    ebegin "${description_stop}"
    start-stop-daemon --stop --quiet --retry=TERM/30/KILL/5 --pidfile ${PIDFILE} --exec ${DAEMON}
    eend $?
}
