# etc/init.d/rsyslog-openrc
#!/sbin/openrc-run
# Copyright (c) 2007-2021 The Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts the rsyslogd enhanced system logging daemon"
command="/usr/sbin/rsyslogd"
command_args="-n ${RSYSLOG_OPTS}"
pidfile="/run/rsyslogd.pid"

# Using supervise-daemon tells OpenRC to monitor this critical daemon
# and automatically handles protecting it from being killed during shutdown.
supervisor=supervise-daemon

depend() {
	# Rsyslog provides the system logger.
	provide logger
	# It should start after basic boot tasks are complete.
	after bootmisc
}
