# etc/init.d/rsyslog-openrc
#!/sbin/openrc-run

description="Starts the rsyslogd enhanced system logging daemon"
command="/usr/sbin/rsyslogd"
command_args="-n"
supervisor=supervise-daemon

depend() {
	provide logger
	after bootmisc
}
