# etc/init.d/rsyslog-openrc
#!/sbin/openrc-run

description="Rsyslog is an enhanced multi-threaded syslogd"

command="/usr/sbin/rsyslogd"
pidfile="/run/rsyslogd.pid"

depend() {
    need mountall
    after bootmisc
}

sendsigs_omit() {
    local omitdir="/run/sendsigs.omit.d"
    mkdir -p "$omitdir"
    ln -sf "${pidfile}" "$omitdir/rsyslog"
}

start() {
    ebegin "Starting Rsyslog"
    sendsigs_omit
    "${command}" &
    echo $! > "${pidfile}"
    eend $?
}

stop() {
    ebegin "Stopping Rsyslog"
    if [ -f "${pidfile}" ]; then
        local pid
        pid=$(cat "${pidfile}")
        if kill "${pid}" 2>/dev/null; then
            rm -f "${pidfile}"
            eend 0
        else
            eend 1
        fi
    else
        eend 1
    fi
}
