# etc/init.d/rsyslog-openrc
#!/sbin/openrc-run

description="Rsyslog is an enhanced multi-threaded syslogd"

command="/usr/sbin/rsyslogd"
pidfile="/run/rsyslogd.pid"

depend() {
    need mountall
    after bootmisc
}

create_xconsole() {
    local xconsole="/dev/xconsole"
    if [ "$(uname -s)" != "Linux" ]; then
        xconsole="/run/xconsole"
        ln -sf "$xconsole" /dev/xconsole
    fi
    if [ ! -e "$xconsole" ]; then
        mknod -m 640 "$xconsole" p
        chown root:adm "$xconsole"
        [ -x /sbin/restorecon ] && /sbin/restorecon "$xconsole"
    fi
}

sendsigs_omit() {
    local omitdir="/run/sendsigs.omit.d"
    mkdir -p "$omitdir"
    ln -sf "${pidfile}" "$omitdir/rsyslog"
}

start() {
    ebegin "Starting Rsyslog"
    create_xconsole
    sendsigs_omit
    echo $! > "${pidfile}"
    eend $?
}

stop() {
    ebegin "Stopping Rsyslog"
    if [ -f "${pidfile}" ]; then
        kill "$(cat "${pidfile}")" && rm -f "${pidfile}"
        eend $?
    else
        eend 1
    fi
}
