# etc/init.d/rsyslog-openrc
#!/sbin/openrc-run

description="Rsyslog is an enhanced multi-threaded syslogd"

command="/usr/sbin/rsyslogd"
pidfile="/run/rsyslogd.pid"

depend() {
    need mountall
    after bootmisc
}

sendsigs_omit() {
    local omitdir="/run/sendsigs.omit.d"
    mkdir -p "$omitdir"
    ln -sf "${pidfile}" "$omitdir/rsyslog"
}

start() {
    ebegin "Starting Rsyslog"
    sendsigs_omit
    echo $! > "${pidfile}"
    eend $?
}

stop() {
    ebegin "Stopping Rsyslog"
    if [ -f "${pidfile}" ]; then
        kill "$(cat "${pidfile}")" && rm -f "${pidfile}"
        eend $?
    else
        eend 1
    fi
}
